# scripts/com.realchoice.tailscaled.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  ⚠ DEPRECATED (2026-05-13, 라운드 12)

  본 LaunchAgent 는 userspace networking 모드 tailscaled 를 가동했습니다.
  라운드 12 에서 `tailscale-app` cask (네이티브 TUN) 로 마이그레이션 — 사용자가
  sudo .pkg 1회 실행 후 메뉴바 앱이 시스템 daemon 운영.

  네이티브 모드 장점:
    - macmini 자체에서도 tailnet URL 도달 가능 (자기 자신 routing)
    - VNC 등 임의 TCP 포트 inbound 라우팅 (옵션 B 의 핵심)
    - 메뉴바 GUI + 자동 갱신

  본 plist 는 **참조용 / 응급 fallback** 으로 보존:
    - 만약 cask 가 어떤 이유로 가동 안 되면 본 plist 로 userspace 모드 복귀.
    - 그 경우 cask 데몬 우선 정지 필요.

  현재 운영:
    - 시스템 daemon : `/Applications/Tailscale.app` (cask)
    - hostname     : realchoice-macmini.tail090015.ts.net
    - serve        : `tailscale serve --bg 8222` → http://127.0.0.1:8222

  이전 사용법 (deprecated):
    cp ~/vault/scripts/com.realchoice.tailscaled.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.realchoice.tailscaled.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.realchoice.tailscaled</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/tailscale/bin/tailscaled</string>
        <string>--tun=userspace-networking</string>
        <string>--socket=/Users/realchoice/Library/Application Support/Tailscale/tailscaled.sock</string>
        <string>--state=/Users/realchoice/Library/Application Support/Tailscale/tailscaled.state</string>
        <string>--statedir=/Users/realchoice/Library/Application Support/Tailscale</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/realchoice/realchoice-ssot/logs/tailscaled.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/realchoice/realchoice-ssot/logs/tailscaled.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/realchoice</string>
        <key>TS_SOCKET</key>
        <string>/Users/realchoice/Library/Application Support/Tailscale/tailscaled.sock</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
