# scripts/com.realchoice.vault-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  vault 백업 LaunchAgent — 매일 03:10
  설치:
    cp ~/vault/scripts/com.realchoice.vault-backup.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.realchoice.vault-backup.plist
  검증:
    launchctl print gui/$(id -u)/com.realchoice.vault-backup
  해제:
    launchctl bootout gui/$(id -u)/com.realchoice.vault-backup
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.realchoice.vault-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>/Users/realchoice/vault/scripts/backup.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>10</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/realchoice/realchoice-ssot/logs/vault-backup.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/realchoice/realchoice-ssot/logs/vault-backup.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/realchoice</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
